# conf_files/^etc^init^restart-vpn.conf
description	"Restart VPN"

start on started vpn-failover
stop on runlevel [!2345]

task

console log

script
    /etc/init.d/openvpn restart
end script

