# conf_files/^etc^init^vpn-event.conf
description	"VPN event hook"

start on vpn-event
stop on runlevel [!2345]

task

console log

exec /usr/bin/python ${BASE}/openvpn-maou/vpn-event.py

