# conf_files/^etc^init^vpn-switch.conf
description	"VPN switch"

start on vpn-event EVENT=vpn_switch
stop on runlevel [!2345]

task

console log

exec ${BASE}/openvpn-maou/switch-vpn.sh $IP

