# conf_files/^etc^init^on-switch-vpn.conf
description	"Random works after switching VPN"

start on vpn-event
stop on runlevel [!2345]

task

console log

exec ${BASE}/openvpn-maou/on-switch-vpn.sh

