# scripts/nat.conf
description	"NAT"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on (runlevel [!2345])

task

chdir /etc/sapikachu/scripts/

console log

script
    ./nat.sh
end script

